# setup/queue.conf
[program:queue]
command=php artisan queue:listen --tries=2 --timeout=600
directory=/var/www/cnp
stdout_logfile=/var/www/cnp/storage/logs/supervisor.log
redirect_stderr=true

