# setup/queue.conf
[program:queue]
command=php artisan queue:listen --tries=2 --timeout=600
directory=/var/www/gbe
stdout_logfile=/var/www/gbe/storage/logs/supervisor.log
redirect_stderr=true

